# deploy/alpha-os.service
[Unit]
Description=Alpha-OS trade runtime (compat wrapper, BTC)
After=network.target

[Service]
Type=simple
User=dev
WorkingDirectory=/home/dev/projects/alpha-os
EnvironmentFile=-/home/dev/.secrets/alpha-os-env
ExecStartPre=/bin/mkdir -p /home/dev/projects/alpha-os/data/BTC/logs /home/dev/projects/alpha-os/data/BTC/metrics
ExecStart=/home/dev/projects/alpha-os/scripts/run_schedule.sh BTC
Restart=on-failure
RestartSec=60
StandardOutput=append:/home/dev/projects/alpha-os/data/BTC/logs/service.log
StandardError=append:/home/dev/projects/alpha-os/data/BTC/logs/service.log

# Graceful shutdown via SIGTERM (PipelineScheduler handles it)
KillSignal=SIGTERM
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
