# deploy/alpha-os-refresh-deployed-alphas@.service
[Unit]
Description=Alpha-OS deployed alphas refresh (%i)
After=network.target

[Service]
Type=oneshot
User=dev
WorkingDirectory=/home/dev/projects/alpha-os
EnvironmentFile=-/home/dev/.secrets/alpha-os-env
ExecStartPre=/bin/mkdir -p /home/dev/projects/alpha-os/data/%i/logs
ExecStart=/bin/bash -lc 'source .venv/bin/activate && exec python3 -m alpha_os_recovery refresh-deployed-alphas --asset %i --config /home/dev/.config/alpha-os/prod.toml'
MemoryHigh=300M
MemoryMax=500M
StandardOutput=append:/home/dev/projects/alpha-os/data/%i/logs/refresh-deployed-alphas.log
StandardError=append:/home/dev/projects/alpha-os/data/%i/logs/refresh-deployed-alphas.log
