# deploy/alpha-os-paper-runtime@.timer
[Unit]
Description=Alpha-OS paper runtime schedule (%i)

[Timer]
OnCalendar=hourly
Persistent=true
Unit=alpha-os-paper-runtime@%i.service

[Install]
WantedBy=timers.target

