# deploy/alpha-os-refresh-universe@.timer
[Unit]
Description=Daily trading universe refresh for %i

[Timer]
OnCalendar=*-*-* 00:40:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
