# deploy/alpha-os-produce-predictions@.service
[Unit]
Description=Alpha-OS hypothesis prediction production (%i)
After=network.target alpha-os-sync-signal-cache@%i.service
Requires=alpha-os-sync-signal-cache@%i.service
OnFailure=alpha-os-runtime-status@%i.service

[Service]
Type=oneshot
WorkingDirectory=/home/dev/projects/alpha-os
EnvironmentFile=-/home/dev/.secrets/alpha-os-env
ExecStartPre=/bin/mkdir -p /home/dev/projects/alpha-os/data/%i/logs
ExecStart=/bin/bash -lc 'source .venv/bin/activate && CONFIG="${ALPHA_OS_CONFIG:-/home/dev/projects/alpha-os/config/dev.toml}" && exec python3 -m alpha_os_recovery produce-predictions --asset %i --strict --config "$CONFIG"'
MemoryHigh=500M
MemoryMax=800M
StandardOutput=append:/home/dev/projects/alpha-os/data/%i/logs/produce-predictions.log
StandardError=append:/home/dev/projects/alpha-os/data/%i/logs/produce-predictions.log
